# overlay.d/99okd/usr/lib/systemd/system/okd-selinux.service
[Service]
Type=oneshot
ExecStart=/usr/sbin/semodule -i /usr/lib/okd/selinux-fixes.cil
RemainAfterExit=yes
[Install]
WantedBy=multi-user.target
